# groq-api-server.service
[Unit]
Description=Groq API Proxy Server (FastAPI)
After=network.target

[Service]
Type=simple
User=gslee
Group=users
WorkingDirectory=/home/gslee/groq-api-server
Environment="PATH=/home/gslee/groq-api-server/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/gslee/groq-api-server/venv/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8021
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
